# nginx/nginx.service
# container-7632c5646ac222f0c5afe52bb60dd78904c32d5c0cc7845dfef249e3f5d5a2f2.service
# autogenerated by Podman 2.2.1
# Sat Dec  5 15:38:21 GMT 2020

[Unit]
Description=Podman container-7632c5646ac222f0c5afe52bb60dd78904c32d5c0cc7845dfef249e3f5d5a2f2.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
ExecStartPre=/bin/rm -f %t/container-7632c5646ac222f0c5afe52bb60dd78904c32d5c0cc7845dfef249e3f5d5a2f2.pid %t/container-7632c5646ac222f0c5afe52bb60dd78904c32d5c0cc7845dfef249e3f5d5a2f2.ctr-id
ExecStart=/usr/bin/podman run --conmon-pidfile %t/container-7632c5646ac222f0c5afe52bb60dd78904c32d5c0cc7845dfef249e3f5d5a2f2.pid --cidfile %t/container-7632c5646ac222f0c5afe52bb60dd78904c32d5c0cc7845dfef249e3f5d5a2f2.ctr-id --cgroups=no-conmon --replace --rm -d --name nginx -p 443:443 -p 80:80 -v /data/nginx/www:/www:Z -v /data/nginx/root/etc/nginx/nginx.conf:/etc/nginx/nginx.conf:Z -v /data/nginx/root/etc/pki/tls:/etc/pki/tls:Z nginx
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/container-7632c5646ac222f0c5afe52bb60dd78904c32d5c0cc7845dfef249e3f5d5a2f2.ctr-id -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/container-7632c5646ac222f0c5afe52bb60dd78904c32d5c0cc7845dfef249e3f5d5a2f2.ctr-id
PIDFile=%t/container-7632c5646ac222f0c5afe52bb60dd78904c32d5c0cc7845dfef249e3f5d5a2f2.pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target
