# unifi/unifi.service
# container-6f7d283158c94b3ec8977c6acb8a25a49472627fe1046b4d379663c62595ae94.service
# autogenerated by Podman 2.2.1

[Unit]
Description=Podman container-6f7d283158c94b3ec8977c6acb8a25a49472627fe1046b4d379663c62595ae94.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
ExecStartPre=/bin/rm -f %t/container-6f7d283158c94b3ec8977c6acb8a25a49472627fe1046b4d379663c62595ae94.pid %t/container-6f7d283158c94b3ec8977c6acb8a25a49472627fe1046b4d379663c62595ae94.ctr-id
ExecStart=/usr/bin/podman run --conmon-pidfile %t/container-6f7d283158c94b3ec8977c6acb8a25a49472627fe1046b4d379663c62595ae94.pid --cidfile %t/container-6f7d283158c94b3ec8977c6acb8a25a49472627fe1046b4d379663c62595ae94.ctr-id --cgroups=no-conmon --replace -d --rm -e PUID=1000 -e PGID=1000 --name unifi -p 3478:3478/udp -p 10001:10001/udp -p 8080:8080 -p 8443:8443 -p 8843:8880 -p 8880:8880 -p 6789:6789 -v /data/unifi/config:/config:Z linuxserver/unifi-controller:latest
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/container-6f7d283158c94b3ec8977c6acb8a25a49472627fe1046b4d379663c62595ae94.ctr-id -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/container-6f7d283158c94b3ec8977c6acb8a25a49472627fe1046b4d379663c62595ae94.ctr-id
PIDFile=%t/container-6f7d283158c94b3ec8977c6acb8a25a49472627fe1046b4d379663c62595ae94.pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target
